# sunnoti.service
[Unit]
Description=Sunrise Notification Bot

[Service]
Type=Simple
ExecStart=/usr/bin/python /your/directory/sunnoti.py
